# config/extra/etc/init.d/media
#!/sbin/openrc-run

description="webdav/sftp/http server over rclone"
command="/usr/bin/rclone"
command_args="rcd --cache-dir /var/cache/rclone --rc-no-auth ~"
command_user="media:media"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need localmount loopback
}

start_post() {
  ${command} rc serve/start type=http addr=:5000 fs=/var/lib/media/public
  ${command} rc serve/start type=http addr=:5010 fs=/var/lib/media user=media pass=${password}
}
